# systemd/monitor-raspi.timer
[Unit]
Description = monitor-raspi Timer

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy = timers.target